# systemd/switchtail-sessions.path
[Unit]
Description=Watch JangLabs submodule config and regenerate SwitchTail kitty session files

[Path]
# .gitmodules changes on submodule add/remove; .git/config on init/deinit.
# The watched tree is the labs WORKSPACE (~/JangLabs), not this repo.
PathModified=%h/JangLabs/.gitmodules
PathModified=%h/JangLabs/.git/config
Unit=switchtail-sessions.service

[Install]
WantedBy=default.target
